# dbprod.pa/etc/cron.hourly/mysql_check_replication.cron
#!/bin/bash

SERVICE='mysqld'

#  Only run the replication check if mysql is running.
if ps ax | grep -v grep | grep $SERVICE > /dev/null
then
  /usr/local/sbin/check_replication.sh
fi
