# webprod.pa/etc/cron.daily/sync_with_dx24n1.sh
#!/bin/sh

/usr/local/sbin/sync_from_dx24n1.sh

