# webalt.pa/etc/cron.daily/sync_webroot_from_webprod.sh
#!/bin/sh

/usr/local/sbin/sync_web_root.sh rsync@webprodpa 'tgannon@gmail.com'
